# nexura-web-maintenance.service
[Unit]
Description=Nexura Digital web + KuntaKinte maintenance check
After=network-online.target nginx.service openclaw-gateway.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/root/nexura/nexuradigital_maintenance_check.sh
Nice=10
IOSchedulingClass=best-effort
